# scripts/systemd/quake3-servers.target
[Unit]
Description=All Quake 3 Servers

[Install]
WantedBy=multi-user.target
