# cmd/trinity/setup/systemd/trinity-digest.timer
[Unit]
Description=Trinity Discord digest, scheduled

[Timer]
# OnCalendar uses the value the install wizard collected. systemd
# parses an enormous range of time specs; "Sun 20:00" is the default.
# Persistent=true lets a missed run (host was down) fire on next boot.
OnCalendar={{schedule}}
Persistent=true

[Install]
WantedBy=timers.target
