# cmd/trinity/setup/systemd/trinity.service
[Unit]
Description=Trinity
# Trinity starts before quake3-servers.target activates the per-server
# instances (which have After=trinity.service). The reverse ordering
# would create a cycle via the per-server units' PartOf= → target →
# After= back to trinity. Log files don't exist until the q3 servers
# spawn, so the tailer's "wait for file" retry covers the launch race.
After=network.target
Before=quake3-servers.target

[Service]
Type=simple
User=quake
Group=quake
ExecStart=/usr/local/bin/trinity serve
Restart=on-failure
RestartSec=5

# Sandboxing. Trinity reads /etc/trinity, writes /var/lib/trinity
# (SQLite + NATS state), and reads /var/log/quake3 for the log tailer.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/trinity
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
